# dev/puppet/files/scanner.conf
[program:scanner]
command=/usr/bin/php /vagrant/cli/command.php scanner
