# dev/puppet/files/scanner.conf
[program:scanner]
command=php /vagrant/cli/command.php scanner
