# dev/puppet/files/web.conf
[program:web]
command=php -S 0.0.0.0:9999 /vagrant/web/index.php -d date.timezone='Europe/Amsterdam'

